# deploy/systemd/wellness-copilot-worker.service
[Unit]
Description=Wellness Copilot WeChat worker
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/wellness-copilot
EnvironmentFile=/opt/wellness-copilot/.env
ExecStart=/usr/bin/python3 scripts/wechat_ilink_worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
